# etc/init.d/shadow
#!/sbin/runscript

depend() {
	need localmount lvm
}

start() {
	ebegin "Activate shadow on the HDD volumes group."
	vgscan
	vgchange -a y
	mount /dev/sdb3 /mnt/shadow
	mount /dev/sdb2 /mnt/shadow/boot
	mount /dev/vghdd1/usr /mnt/shadow/usr
	mount /dev/vghdd1/var /mnt/shadow/var
	mount /dev/vghdd1/opt /mnt/shadow/opt
	mount /dev/vghdd1/tmp /mnt/shadow/tmp
	mount /dev/vghdd1/home /mnt/shadow/home
	mount /dev/vghdd1/usrlocal /mnt/shadow/usr/local
	mount /dev/vghdd1/portage /mnt/shadow/usr/portage
	mount /dev/vghdd1/distfiles /mnt/shadow/usr/portage/distfiles
	mount /dev/vghdd1/vartmp /mnt/shadow/var/tmp
	mount -t proc none /mnt/shadow/proc
	mount --rbind /dev /mnt/shadow/dev
	cp -L /etc/resolv.conf /mnt/shadow/etc/
	eend $?
}

stop() {
	ebegin "Deactivate shadow on the HDD volumes group."
	umount /mnt/shadow/{usr/local,var/tmp,tmp,opt,home,dev/pts,dev/shm,proc,boot,usr/portage/distfiles}
	umount /mnt/shadow/{usr/portage,var,dev}
	umount /mnt/shadow/usr
	umount /mnt/shadow
	eend $?
}
